# launchd/com.mms.docker-gc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Docker のゴミ掃除を日次で回す（2026-07-23 追加）

  ★なぜ要るか
    ビルドのたびに旧イメージとビルドキャッシュが残り、放置すると数十GBになる。
    2026-07-23 にディスクが99%まで埋まり、Postgres がチェックポイントを書けず
    クラッシュループした。DBが落ちると web も worker も同時に止まり、
    ジョブの成否では気づけない（ジョブ自体が動けないので記録も残らない）。

  ★毎日 04:30。日次ジョブ（05:00〜）より前に済ませて、
    ジョブが動く時点で空きがある状態にする。

  インストール:
    cp launchd/com.mms.docker-gc.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.mms.docker-gc.plist

  アンインストール:
    launchctl unload ~/Library/LaunchAgents/com.mms.docker-gc.plist
    rm ~/Library/LaunchAgents/com.mms.docker-gc.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mms.docker-gc</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ishiimasataka/システム開発/Next/media-management-system/scripts/docker-gc.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>4</integer>
    <key>Minute</key><integer>30</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/ishiimasataka/システム開発/Next/media-management-system/logs/docker-gc.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ishiimasataka/システム開発/Next/media-management-system/logs/docker-gc.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/ishiimasataka/システム開発/Next/media-management-system</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/Applications/Docker.app/Contents/Resources/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TZ</key>
    <string>Asia/Tokyo</string>
  </dict>
</dict>
</plist>
